# src/crontab
0 3 * * * python /usr/app/src/main.py >> /logfile